# src/third_party/mozjs-45/js/src/devtools/rootAnalysis/tools/clobberer/cleanup_clobberer.cron
MAILTO=catlee@mozilla.com
@hourly root CLEANUP_PASSWORD=sekrits /usr/bin/php /var/www/html/build/clobberer/cleanup.php
